# systemd/bot-slop.service
[Unit]
Description=AI Workflow Slop Bot Daemon - Code Quality Monitor
After=network.target dbus.socket
Wants=dbus.socket
# Max 5 restarts in 5 minutes, then give up
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
# Type=notify allows daemon to signal readiness and supports watchdog
Type=notify
NotifyAccess=main
# Watchdog: systemd restarts if no ping received within 15s
WatchdogSec=15

WorkingDirectory=/home/daoneill/src/redhat-ai-workflow
# Environment variables for Claude (Vertex AI)
Environment=CLAUDE_CODE_USE_VERTEX=1
Environment=ANTHROPIC_VERTEX_PROJECT_ID=itpc-gcp-ansible-pe-eng-claude
Environment=ANTHROPIC_VERTEX_REGION=us-east5
# Display environment for notifications (X11/Wayland)
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
# XDG runtime for D-Bus/notifications
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Use the virtual environment's python with D-Bus enabled
ExecStart=/home/daoneill/src/redhat-ai-workflow/.venv/bin/python3 -m services.slop --dbus

# Restart policy
Restart=on-failure
RestartSec=5

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed

# Logs go to journald (view with: journalctl --user -u bot-slop -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
