# systemd/bot-meet.service
[Unit]
Description=AI Workflow Google Meet Bot Daemon
After=network.target dbus.socket
Wants=dbus.socket
# Max 5 restarts in 5 minutes, then give up
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/daoneill/src/redhat-ai-workflow
# Environment variables for Claude (Vertex AI)
Environment=CLAUDE_CODE_USE_VERTEX=1
Environment=ANTHROPIC_VERTEX_PROJECT_ID=itpc-gcp-ansible-pe-eng-claude
Environment=ANTHROPIC_VERTEX_REGION=us-east5
# Environment variables for Gemini CLI (Vertex AI)
Environment=GOOGLE_CLOUD_PROJECT=itpc-gcp-ansible-pe-eng-claude
Environment=GOOGLE_CLOUD_LOCATION=global
# Display environment for browser (X11/Wayland)
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
# XDG runtime for PipeWire/PulseAudio access
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Use the virtual environment's python with D-Bus enabled
ExecStart=/home/daoneill/src/redhat-ai-workflow/.venv/bin/python3 scripts/meet_daemon.py --dbus

# Restart policy
Restart=on-failure
RestartSec=30

# Give the daemon 45 seconds to shut down gracefully before SIGKILL
# This allows time to leave meetings and close browsers
TimeoutStopSec=45
# Send SIGTERM first, then SIGKILL after TimeoutStopSec
KillMode=mixed

# Logs go to journald (view with: journalctl --user -u bot-meet -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
