# systemd/bot-cron.service
[Unit]
Description=AI Workflow Cron Scheduler Daemon
After=network.target dbus.socket
Wants=dbus.socket
# Max 5 restarts in 5 minutes, then give up
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/daoneill/src/redhat-ai-workflow
# Environment variables for Claude (Vertex AI)
Environment=CLAUDE_CODE_USE_VERTEX=1
Environment=ANTHROPIC_VERTEX_PROJECT_ID=itpc-gcp-ansible-pe-eng-claude
Environment=ANTHROPIC_VERTEX_REGION=us-east5
# Display environment for notifications (X11/Wayland)
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
# XDG runtime for D-Bus/notifications
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Use the virtual environment's python with D-Bus enabled
ExecStart=/home/daoneill/src/redhat-ai-workflow/.venv/bin/python3 scripts/cron_daemon.py --dbus

# Restart policy
Restart=on-failure
RestartSec=30

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed

# Logs go to journald (view with: journalctl --user -u bot-cron -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
