# systemd/extension-watcher.service
[Unit]
Description=AI Workflow VSCode Extension Watcher
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/daoneill/src/redhat-ai-workflow
# Display environment for notifications (X11/Wayland)
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
# XDG runtime for D-Bus/notifications
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/home/daoneill/src/redhat-ai-workflow/.venv/bin/python3 -m services.extension_watcher --notify
Restart=always
RestartSec=10
# Logs go to journald (view with: journalctl --user -u extension-watcher -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
